# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/php/php-fpm.log
stderr_logfile=/var/log/php/php-fpm-error.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off; pid /tmp/nginx.pid;"
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log

[program:queue-worker]
command=/usr/local/bin/php /var/www/html/artisan queue:work --sleep=1 --tries=3 --max-time=3600 --timeout=90
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/queue-worker.log
stderr_logfile=/var/log/supervisor/queue-worker-error.log

[program:scheduler]
command=/usr/local/bin/php /var/www/html/artisan schedule:work
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/scheduler.log
stderr_logfile=/var/log/supervisor/scheduler-error.log
